# com.financehub.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.financehub.scheduler</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/home/Documents/Info Site/mubasher-deep-extract/scheduler.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/home/Documents/Info Site/mubasher-deep-extract</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/home/Documents/Info Site/mubasher-deep-extract/logs/scheduler_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/home/Documents/Info Site/mubasher-deep-extract/logs/scheduler_stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
